# scripts/thermax-browser.service
[Unit]
Description=Thermax HMI Browser Launcher
After=thermax-hmi.service network.target
Requires=thermax-hmi.service

[Service]
Type=oneshot
User=pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
ExecStart=/bin/bash /home/pi/ThermaxDashboard/local_server/scripts/open_browser.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target

